# com.pixelancer.sheets-watcher.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.pixelancer.sheets-watcher</string>

    <key>ProgramArguments</key>
    <array>
        <!-- Replace with the full path to your run.sh file -->
        <string>/Users/YOUR_USERNAME/path/to/sheets-folder-watcher/run.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <!-- Replace with the full path to your sheets-folder-watcher directory -->
    <string>/Users/YOUR_USERNAME/path/to/sheets-folder-watcher</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <!-- Replace with the full path where you want the log file saved -->
    <string>/Users/YOUR_USERNAME/path/to/sheets-folder-watcher/launchd.log</string>

    <key>StandardErrorPath</key>
    <!-- Replace with the full path where you want the error log saved -->
    <string>/Users/YOUR_USERNAME/path/to/sheets-folder-watcher/launchd-error.log</string>
</dict>
</plist>
